# config/service_script.conf
[supervisord]  ## This is the main process for the Supervisor
nodaemon=true  ## This setting is to specify that we are not running in daemon mode

## agora_service
[program:agora_service]
command=agora -c %(ENV_CONFIG)
autorestart=true
stderr_logfile=/dev/stdout
stderr_logfile_maxbytes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

## stoa_service
[program:stoa_service]
command=npm start
autostart=true
autorestart=true
stderr_logfile =/dev/stdout
stderr_logfile_maxbytes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
